# traitement_texte/profile.d/jcenv.sh
alias jcenv='. /usr/local/libexec/jcenv'
