# traitement_texte/profile.d/path.sh
# {{ ansible_managed }}

[ $UID -eq 0 ] && PATH=$PATH:/local/sbin
PATH=$PATH:/local/bin

