# traitement_texte/profile.d/SU.sh
# Don't edit, this file was generated by ansible

alias SU="echo 'tips: next time, try \"sudo -s\" instead'; sudo -s"
